# backend/cron/clear_cache.cron
0 3 * * * cd /app && python scripts/clear_fastf1_cache.py

